# docker/worker/crontab
*/30 * * * * /usr/local/bin/php /var/www/bin/console app:parse >> /var/www/docker/log/cron.log 2>&1
