# docker/worker/conf.d/cron.conf
[program:cron]
command=/usr/sbin/cron -f
autostart=true
autorestart=true
stderr_logfile=/var/www/var/log/cron.err.log
stdout_logfile=/var/www/var/log/cron.out.log
